# config/hello.conf
[program:hello]
command=/home/admin/venv/bin/gunicorn HelloRekruro.wsgi:application -c /home/admin/HelloRekruro/config/gunicorn.conf.py
directory=/home/admin/HelloRekruro
user=admin
autorestart=true
redirect_stderr=true
stdout_logfile = /home/admin/HelloRekruro/logs/debug.log